# config/supervisor/webservice.conf
[program:webservice]
command=npm start
user=pi
stderr_logfile = /var/log/supervisor/webservice-err.log
stdout_logfile = /var/log/supervisor/webservice-stdout.log
directory=/home/pi/RPiTC/webserver

